# Assignment-2/tomcat.service
[Unit]
Description=Apache Tomcat Web Application Container
After=network.targetd

[Service]
Type=forking